# config/kspw-web/cron/renew-certs.cron
##############################
#    Renew Certs Cronjobs    #
##############################

# At every minute
* * * * * echo "Heartbeat $(date)" >> /var/log/cron.log 2>&1

# At every 2nd minute
# */2 * * * * certbot-auto renew >> /var/log/cron.log 2>&1

# At minute 0 past hour 0 and 12
# 0 0,12 * * * certbot-auto renew >> /var/log/cron.log 2>&1

#
