# config/kspw-db/cron/upload-cloud.cron
###############################
#    Cloud Upload Cronjobs    #
###############################

# At every 2nd minute
# */2 * * * * bash /home/bin/upload_cloud.sh >> /var/log/cron.log 2>&1

# Daily at 23:59
59 23 * * * bash /home/bin/upload_cloud.sh >> /var/log/cron.log 2>&1

# Weekly on Sunday at 23:59
# 59 23 * * 0 bash /home/bin/upload_cloud.sh >> /var/log/cron.log 2>&1

# Monthly on 1st at 23:59
# 59 23 1 * * bash /home/bin/upload_cloud.sh >> /var/log/cron.log 2>&1

#
